# gmail-watcher.service
[Unit]
Description=Gmail Watcher - AI Employee
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=/mnt/c/Users/Admin/Documents/GitHub/Hackathon_0_Personal_AI_Employee_FTE
ExecStart=/mnt/c/Users/Admin/Documents/GitHub/Hackathon_0_Personal_AI_Employee_FTE/venv/bin/python watchers/gmail_watcher.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/gmail_watcher.log
StandardError=append:/var/log/gmail_watcher.log

[Install]
WantedBy=multi-user.target
